# .config/systemd/user/backup.service
[Unit]
Description=Backup Service

[Service]
Type=simple
ExecStart= /home/karsten/.libshell/backup create > /dev/null
