# .config/systemd/user/monitor.timer
[Unit]
Description=Timer for monitoring service

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target
