# .config/autostart/ptyxis.desktop
[Desktop Entry]
Type=Application
Exec=flatpak run app.devsuite.Ptyxis
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=Ptyxis
